# debian/pymc-repeater.service
[Unit]
Description=PyMC Repeater Daemon
After=network.target

[Service]
Type=simple
User=pymc-repeater
Group=pymc-repeater
WorkingDirectory=/etc/pymc-repeater
ExecStart=/usr/bin/pymc-repeater
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
